# aws/packer/Web/scripts/ROOT/etc/supervisor/conf.d/check_shared.conf
[program:check_shared]
command=/root/check.sh -t 10 -m /var/www/flipit.com/shared
user=root
