# aws/packer/Web/scripts/ROOT/etc/supervisor/conf.d/check_backup.conf
[program:check_backup]
command=/root/check.sh -t 60 -m /var/www/flipit.com/backup
user=root
