# aws/packer/Web/scripts/ROOT/etc/supervisor/conf.d/check_quota.conf
[program:check_quota]
command=/root/quota.sh -t 30 -m /mnt -q 70
user=root
