# aws/packer/Web/scripts/ROOT/etc/supervisor/conf.d/check_health.conf
[program:check_health]
command=/root/check_health.sh -t 60 
user=root
