# aws/packer/Web/scripts/ROOT/etc/supervisor/conf.d/check_language.conf
[program:check_language]
command=/root/check.sh -t 60 -m /var/www/flipit.com/language
user=root
